# docker/dhcpd/config/supervisord_app.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
loglevel=debug
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:dhcpd]
command=/opt/cnaas/dhcpd.sh
autorestart=true
